# system-files/usr/lib/systemd/system/bootc-upgrade.timer
[Unit]
Description=Periodic bootc image upgrade
Documentation=man:bootc(8)

[Timer]
OnBootSec=10min
OnUnitInactiveSec=1d
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
